# packaging/systemd/user/hyos-docd.service
[Unit]
Description=HyOS Document service — document analysis and drafting
Documentation=https://hyos.tech/docs/services/docd
After=dbus.socket

[Service]
Type=dbus
BusName=org.hyos.Documents1
ExecStart=%h/.local/bin/hyos-docd

# --- Hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
# AF_INET needed for Ollama on localhost
RestrictAddressFamilies=AF_UNIX AF_INET
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RemoveIPC=yes

# Only Ollama on localhost
IPAddressAllow=127.0.0.1/8 ::1/128
IPAddressDeny=any

# No write paths — this service writes nothing
ReadWritePaths=

[Install]
WantedBy=default.target
