# packaging/systemd/user/hyos-indexerd.service
[Unit]
Description=HyOS Indexer — local document content indexer
Documentation=https://hyos.tech/docs/services/indexerd
After=dbus.socket

[Service]
Type=dbus
BusName=org.hyos.Indexer1
ExecStart=%h/.local/bin/hyos-indexerd

# --- Hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
# ProtectHome=read-only gives the service read access to home; path allow-listing
# is enforced in Python using XDG user-dirs (locale-aware).
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
# No network whatsoever — indexer does not call Ollama
RestrictAddressFamilies=AF_UNIX
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RemoveIPC=yes

# Index storage (the only writable path)
ReadWritePaths=%h/.local/share/hyos/index

[Install]
WantedBy=default.target
